# root/etc/crontabs/abc
# file expiration job
0 0 * * * find /data -type f -mtime +7 -delete \;
